# initramfs/etc/profile.d/move_ssh_auth.sh
#!/bin/sh

mkdir ~/.ssh
mv /etc/dropbear/authorized_keys ~/.ssh
